# config/systemd/user/resethdmiout.service
[Service]
Type=oneshot
RemainAfterExit=true
ExecStop=rm /home/drew/.asoundrc

[Install]
WantedBy=default.target
